# _archive/deploy/deploy/diagnosis-monitoring.service
[Unit]
Description=Diagnosis System Monitoring Daemon
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/PythonProject
Environment="PATH=/usr/bin:/usr/local/bin"
EnvironmentFile=/path/to/PythonProject/.env
ExecStart=/usr/bin/python3 monitoring_daemon.py
Restart=always
RestartSec=10
StandardOutput=append:/path/to/PythonProject/logs/monitoring.log
StandardError=append:/path/to/PythonProject/logs/monitoring.log

# 安全设置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
